# systemd/kiro-auth-refresh.service
[Unit]
Description=Kiro auth token refresh — syncs SSO token via kirolink
Documentation=https://github.com/alexandeism/kirolink
After=network.target

[Service]
Type=oneshot
ExecStart=%h/Documents/proxy/kirolink/kirolink refresh
WorkingDirectory=%h/Documents/proxy/kirolink
